# bologan-bot.service
[Unit]
Description=Victor Bologan Telegram Bot
After=network.target

[Service]
Type=simple
User=YOUR_LINUX_USERNAME
WorkingDirectory=/home/YOUR_LINUX_USERNAME/telegram-bot
EnvironmentFile=/home/YOUR_LINUX_USERNAME/telegram-bot/.env
ExecStart=/home/YOUR_LINUX_USERNAME/telegram-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
